# rootfs/etc/services.d/postfix/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

set -eo pipefail

# Wait for clamd & clamav-milter
if [ "${ENABLE_CLAMAV}" = "true" ]; then
    s6-svwait -U /run/s6/services/clamd /run/s6/services/clamav-milter
fi

# Wait for opendkim
if [ "${ENABLE_OPENDKIM}" = "true" ]; then
    s6-svwait -U /run/s6/services/opendkim
fi

# Wait for postgrey
if [ "${ENABLE_POSTGREY}" = "true" ]; then
    s6-svwait -U /run/s6/services/postgrey
fi

# Wait for syslog
s6-svwait -U /run/s6/services/syslogd

# Start
s6-notifyoncheck \
        -s 1000 \
        -t 1000 \
        -n 0 \
        -c "/usr/local/bin/check_postfix.sh" \
        /usr/sbin/postfix start-fg \
        2>&1 | awk -W Interactive '{print "[postfix] " $0}'
