# rootfs/etc/services.d/clamav-milter/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -eo pipefail

# # Do we run clamav
if [ "${ENABLE_CLAMAV}" = "true" ]; then

    # Wait for clamav & syslog to start
    s6-svwait -U /run/s6/services/clamd /run/s6/services/syslog

    s6-notifyoncheck \
        -s 1000 \
        -t 1000 \
        -n 0 \
        -c "/usr/local/bin/check_clamd.sh" \
        clamav-milter \
            --config-file="${CLAMAV_MILTERCONF_FILE}" \
            2>&1 | awk -W Interactive '{print "[clamav-milter] " $0}'
else
    sleep 86400
    exit 0
fi
